# aiops.service
[Unit]
Description=AIOps Coralogix Monitor
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/aiops-venv
EnvironmentFile=/etc/aiops.env
ExecStart=/home/ubuntu/aiops-venv/bin/python /home/ubuntu/aiops-venv/aiops_webhook.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
